# docs/examples/nut-server.service.d/configure-nut-users.conf.in
[Service]
ExecStartPre=@libexecdir@/@PACKAGE@/setup-nut-users-ExecStartPre.sh
