# docs/examples/nut-server.service.d/configure-nut-users.conf.in
[Service]
ExecStartPre=@datadir@/@PACKAGE@/scripts/setup-nut-users-ExecStartPre.sh
