# systemd/bios-agent-inventory.service.in
[Unit]
Description=Inventory agent for 42ity project
After=malamute.service bios-db-init.service
Requires=malamute.service bios-db-init.service
PartOf=bios.target

[Service]
Type=simple
User=bios
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios-db-rw
ExecStart=@libexecdir@/@PACKAGE@/agent-inventory
Restart=always

[Install]
WantedBy=bios.target
