# systemd/bios.service
[Unit]
Description=Service to toggle the BIOS pseudo-target and whoever wants it
Requires=multi-user.target network.target
After=multi-user.target network.target
Conflicts=rescue.target shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=/bin/sh -c '/bin/systemctl stop -- $(/bin/systemctl show -p Wants bios.service | cut -d= -f2 | sed "s/[ \t]*\(system\.slice\)//")'

[Install]
WantedBy=multi-user.target
