# systemd/bios-agent-cm.service.in
[Unit]
Description=Computation Module agent for BIOS project
After=malamute.service bios-db-init.service
Requires=malamute.service bios-db-init.service

[Service]
Type=simple
User=nobody
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/sysconfig/bios
EnvironmentFile=-/etc/default/bios__%n.conf
# Suggested custom config:
# - BIOS_LOG_LEVEL=warning set a quieter bios log level
#
EnvironmentFile=-/etc/sysconfig/bios__%n.conf
EnvironmentFile=-/etc/default/bios-db-ro
ExecStart=@libexecdir@/@PACKAGE@/agent-cm

[Install]
WantedBy=bios.target
