# systemd/biospath-messages.service.in
[Unit]
Description=BIOS path to fix persmission of /var/log/messages
Conflicts=shutdown.target

[Service]
Type=simple
ExecStart=/bin/chown root:adm /var/log/messages

[Install]
WantedBy=multi-user.target
