# systemd/bios-db-init.service.in
[Unit]
Description=Ensure database schema exists for BIOS services
After=mysql.service
Requires=mysql.service
Conflicts=shutdown.target

[Service]
# it is expected that the process has to exit before systemd starts follow-up units
Type=oneshot
# the service shall be considered active even when all its processes exited
RemainAfterExit=yes
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/sysconfig/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/sysconfig/bios__%n.conf
ExecStart=@libexecdir@/@PACKAGE@/db-init

[Install]
WantedBy=bios.service
