# systemd/bios-networking.service.in
[Unit]
Description=Async restart of networking

[Service]
Type=simple
SyslogIdentifier=%N
ExecStart=@libexecdir@/@PACKAGE@/bios-networking
