# systemd/biostimer-compress-logs.service.in
[Unit]
Description=42ity-Timer service that compresses older log files
Conflicts=rescue.target shutdown.target poweroff.target halt.target reboot.target emergency.target

[Service]
Type=simple
ExecStart=@libexecdir@/@PACKAGE@/compress-logs

[Install]
WantedBy=multi-user.target
