# systemd/bios-fake-th.service.in
[Unit]
Description=Generates fake temperature and humidity readings
After=bios-agent-dbstore.service
Requires=bios-agent-dbstore.service
PartOf=bios.target

[Service]
Type=simple
User=nobody
ExecStart=@libexecdir@/@PACKAGE@/fake-th

[Install]
WantedBy=bios.target
