# systemd/bios-agent-alert.service.in
[Unit]
Description=Alert agent for BIOS project
After=malamute.service
Requires=malamute.service

[Service]
Type=simple
User=nobody
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/sysconfig/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/sysconfig/bios__%n.conf
EnvironmentFile=-/etc/default/bios-db-ro
ExecStart=@libexecdir@/@PACKAGE@/agent-alert

[Install]
WantedBy=bios.target
