# systemd/biospath-messages.path
[Unit]
Description=BIOS path to fix persmission of /var/log/messages
PartOf=multi-user.target

[Path]
PathExists=/var/log/messages

[Install]
WantedBy=multi-user.target
