# systemd/biostimer-warranty-metric.service.in
[Unit]
Description=42ity-Timer service to send metrics about warranty date
After=malamute.service fty-db-init.service
Requires=malamute.service fty-db-init.service
PartOf=bios.target

[Service]
Type=simple
User=bios
# TODO: Is it really a timer or an ordinary service? What is it intended to be?
#Restart=always
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
Environment="prefix=@prefix@"
EnvironmentFile=-@sysconfdir@/default/bios-db-ro
PermissionsStartOnly=true
ExecStart=@libexecdir@/@PACKAGE@/warranty-metric

[Install]
WantedBy=bios.target
