# systemd/bios-agent-alert-smtp.service.in
[Unit]
Description=Notification agent for BIOS project
After=malamute.service
Requires=malamute.service

[Service]
Type=simple
User=nobody
#
# notification environment variables are
# - BIOS_EMAIL_FROM email address of the sender [bios@eaton.com]
# - BIOS_EMAIL_TO comma separated email addresses of recipients [no default]
# - BIOS_EMAIL_BODY_START body of the alert email [Alert ${rulename} started.]
# - BIOS_EMAIL_SUBJECT_START subject of the alert email [Alert P${priority} ${rulename} is active!]
# - BIOS_EMAIL_BODY_END body of the end-alert email [Alert ${rulename} ended.]
# - BIOS_EMAIL_SUBJECT_END subject of the end-alert email [Alert P${priority} ${rulename} ended.]
# - BIOS_SMTP_SERVER SMTP server for relaying emails [no default]
#
# You can use ${rulename} and ${priority} variables in texts (case sensitive, curly brackets are mandatory)
# Put those variables into config file in /etc/sysconfig
#
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/sysconfig/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/sysconfig/bios__%n.conf
ExecStart=@libexecdir@/@PACKAGE@/agent-alert-smtp

[Install]
WantedBy=multi-user.target
