# systemd/bios-agent-inventory.service.in
[Unit]
Description=Testing agent for CI purposes only

[Service]
Type=simple
User=nobody
### the service shall be considered active even when all its processes exited
RemainAfterExit=yes
ExecStart=/bin/cat

[Install]
WantedBy=bios.target
