# systemd/bios-agent-inventory.service.in
[Unit]
Description=Testing agent for CI purposes only

[Service]
Type=simple
User=bios
ExecStart=/bin/cat

[Install]
WantedBy=bios.target
