# systemd/biospath-compress-logs.service.in
[Unit]
Description=BIOS-Path service that asynchronously compresses older log files
Conflicts=shutdown.target

[Service]
Type=simple
ExecStart=@libexecdir@/@PACKAGE@/compress-logs

[Install]
WantedBy=multi-user.target
