# systemd/fty-license-accepted.path.in
# This unit fires whenever the monitored pathname is written to, with its
# action being to try enabling the target Unit. The latter can also start
# by itself during subsequent boots, thanks to the file already existing.
# The target Unit by itself also has a Condition on this file to be not empty.

[Unit]
Description=Tracker that the 42ity license is accepted
PartOf=fty-license-accepted.target

[Path]
PathModified=@ftydatadir@/fty-eula/license
PathExists=@ftydatadir@/fty-eula/license
Unit=fty-license-accepted.service

[Install]
RequiredBy=fty-license-accepted.target
