# systemd/fty-proxy.service.in
[Unit]
Description=Track that the proxy config was changed, for dependent units to react
After=local-fs.target
Requires=local-fs.target
PartOf=bios.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/bin/ls -la @sysconfdir@/default/fty-proxy

[Install]
WantedBy=bios.target
