# systemd/fty-proxy.service.in
[Unit]
Description=Track that the proxy config was changed, for dependent units to react
After=local-fs.target
Requires=local-fs.target
PartOf=bios.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/bin/ls -la @sysconfdir@/default/fty-proxy
ExecStartPost=/bin/dash -c "/bin/systemctl try-restart --no-block -- $(/bin/systemctl show -p WantedBy -p RequiredBy -p BoundBy fty-proxy.service | cut -d= -f2 | sort | uniq | tr '\n' ' ')"

[Install]
WantedBy=bios.target
