# systemd/fty-hostname-setup.service.in
[Unit]
Description=Prepare standard hostname for a 42ity deployment, if none is set statically
After=local-fs.target
Requires=local-fs.target
# Hopefully, being WantedBy+Before multi-user.target allows the naming to
# be set up before the bulk of 42ity product services even begins starting
Before=network-online.target multi-user.target

[Service]
Type=oneshot
# the service shall be considered active even when all its processes exited
RemainAfterExit=yes
Restart=no
User=root
Group=root
# Generate a name unless one is set, do not save to /etc/hostname yet (so DHCP may override)
ExecStart=@libexecdir@/@PACKAGE@/fty-hostname-setup "-" false

[Install]
WantedBy=network-online.target multi-user.target
