# systemd/fty-db-init.service.in
[Unit]
Description=Initialize database schema for 42ity services
After=mysql.service
Requires=mysql.service
Conflicts=shutdown.target
#do not run if database was not yet initialized
ConditionPathExists=/var/lib/fty/sql/mysql/bios-db-rw

[Service]
# it is expected that the process has to exit before systemd starts follow-up units
Type=oneshot
User=root
# the service shall be considered active even when all its processes exited
RemainAfterExit=yes
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
Environment="prefix=@prefix@"
ExecStart=@libexecdir@/@PACKAGE@/db-init

[Install]
WantedBy=bios.target
Alias=bios-db-init.service
