# systemd/fty-proxy.service.in
[Unit]
Description=Prepare common run-time environment variables for 42ity services
After=local-fs.target
Requires=local-fs.target
PartOf=bios.target

[Service]
Type=oneshot
# the service shall be considered active even when all its processes exited
RemainAfterExit=yes
Restart=no
User=root
Group=root
ExecStart=/bin/ls -la @sysconfdir@/default/fty-proxy

[Install]
WantedBy=bios.target
