# systemd/fty-db-init.service.in
[Unit]
Description=Initialize or update database schema for 42ity services
After=fty-db-engine.service
Requires=fty-db-engine.service
BindsTo=fty-db-engine.service
Conflicts=shutdown.target recovery.target
PartOf=bios.target

[Service]
Type=simple
User=root
# the service shall be considered active even when all its processes exited
RemainAfterExit=yes
Restart=always
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/usr/share/fty/etc/default/fty
EnvironmentFile=-/usr/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/fty
EnvironmentFile=-/etc/default/fty__%n.conf
Environment="prefix=/usr"
ExecStart=/usr/libexec/bios/db-init
ExecStopPost=/bin/dash -c "/bin/systemctl stop -- $(/bin/systemctl show -p WantedBy -p RequiredBy -p BoundTo fty-db-init.service | cut -d= -f2 | tr ' ' '\\n' | egrep -v '^(bios|fty)\.(service|target)$')"

[Install]
WantedBy=bios.target
Alias=bios-db-init.service
