# systemd/fty-tntnet@.service.in
[Unit]
Description=Tntnet web server instance for 42ity using @sysconfdir@/tntnet/%i.xml
After=network.target saslauthd.service malamute.service fty-envvars.service
Requires=saslauthd.service malamute.service fty-envvars.service
PartOf=bios.target

[Service]
Type=simple
Restart=always
User=root
Group=root
UMask=0000
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
# Note: for webserver the DB credentials file is optional,
# since acceptance of license via REST API starts the database
EnvironmentFile=-@sysconfdir@/default/bios-db-rw
Environment='SYSTEMD_UNIT_FULLNAME=%n'
PrivateTmp=true
ExecStartPre=@datadir@/@PACKAGE@/scripts/tntnet-ExecStartPre.sh %i
EnvironmentFile=-/run/tntnet-%i.env
EnvironmentFile=-/run/fty-envvars.env
ExecStart=/usr/bin/tntnet -c @sysconfdir@/tntnet/%i.xml
RestartSec=20

[Install]
WantedBy=bios.target
#Alias=tntnet@.service
#Alias=tntnet@%i.service
