# systemd/ifplug-dhcp-autoconf.service.in
[Unit]
Description=Preconfigure ifplugd with currently tracked interfaces and yank them from default networking
DefaultDependencies=no
After=mountkernfs.service local-fs.target urandom.service
Before=ifplugd.service networking.service
Before=network.target
Conflicts=shutdown.target poweroff.target halt.target reboot.target
PropagatesReloadTo=ifplugd.service
ReloadPropagatedFrom=networking.service
# This unit must be alone on its line (can be edited during OS image generation)
Before=ifupdown-wait-online.service

[Service]
Type=simple
### the service shall be considered active even when all its processes exited
RemainAfterExit=yes
User=root
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
Environment="prefix=@prefix@"
ExecStart=@datadir@/@PACKAGE@/scripts/ifplug-dhcp-autoconf.sh

[Install]
WantedBy=network.target
WantedBy=ifplugd.service networking.service
# This unit must be alone on its line (can be edited during OS image generation)
WantedBy=ifupdown-wait-online.service
