# systemd/ifplug-dhcp-autoconf.service.in
[Unit]
Description=Preconfigure ifplugd with currently tracked interfaces and yank them from default networking
DefaultDependencies=no
After=mountkernfs.service local-fs.target urandom.service
Before=ifplugd.service networking.service
Conflicts=shutdown.target
PropagatesReloadTo=ifplugd.service
ReloadPropagatedFrom=networking.service

[Service]
Type=simple
### the service shall be considered active even when all its processes exited
RemainAfterExit=yes
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/sysconfig/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/sysconfig/bios__%n.conf
ExecStart=@datadir@/@PACKAGE@/scripts/ifplug-dhcp-autoconf.sh

[Install]
WantedBy=bios.service
