# systemd/ipc-meta-setup.service.in
[Unit]
Description = Initial IPC setup on target system
DefaultDependencies=no
After = local-fs.target
Requires = local-fs.target
Before = bios.target fty-envvars.service fty-license-accepted.service bios.service malamute.service sshd.service ssh.service ssh.socket systemd-logind.service network.target systemd-tmpfiles-setup.service systemd-tmpfiles-clean.service tntnet@bios.service fty-tntnet@bios.service fty-db-engine.service mysql.service mysqld.service mariadb.service saslauthd.service 
#Requires = bios.target bios.service

[Service]
Type = oneshot
User = root
ExecStart = @prefix@/share/fty/setup/ipc-meta-setup.sh
RemainAfterExit = yes
PrivateTmp = yes

[Install]
RequiredBy=network.target systemd-tmpfiles-setup.service
