# systemd/legacy-asset.service.in
[Unit]
Description=Legacy asset for 42ity project
After=malamute.service bios-db-init.service tntnet@bios.service
Requires=malamute.service bios-db-init.service
PartOf=bios.target

[Service]
Type=simple
User=bios
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios-db-rw
ExecStart=@libexecdir@/@PACKAGE@/legacy-asset
Restart=always

[Install]
WantedBy=bios.target
