# tools/biostimer-graphs-prefetch.service.in
[Unit]
Description=BIOS-Timer service that prefetches average computations for quicker graphs of ALL THERE IS (implementation)
# This service only relies on REST API server directly
After=tntnet@bios.service bios-db-init.service
Requires=tntnet@bios.service bios-db-init.service

[Service]
Type=simple
User=bios
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/sysconfig/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/sysconfig/bios__%n.conf
EnvironmentFile=-/etc/default/bios-db-ro
ExecStart=@datadir@/@PACKAGE@/scripts/biostimer-graphs-prefetch.sh

[Install]
WantedBy=multi-user.target
