# tools/bios-agent-cm.service.in
[Unit]
Description=Computation Module agent for BIOS project
After=malamute.service mysql.service
Requires=malamute.service mysql.service

[Service]
Type=simple
#
# - BIOS_LOG_LEVEL set bios log level
#
EnvironmentFile=-/etc/sysconfig/bios
EnvironmentFile=-/etc/default/bios
ExecStart=@libexecdir@/bios/agent-cm

[Install]
WantedBy=multi-user.target
