# tools/bios-agent-dbstore.service.in
[Unit]
Description=Database layer for BIOS
After=malamute.service bios-db-init.service
Requires=malamute.service bios-db-init.service

[Service]
Type=simple
ExecStart=@libexecdir@/bios/agent-dbstore

[Install]
WantedBy=multi-user.target
