# tools/bios-db-init.service.in
[Unit]
Description=Ensure database schema exists for BIOS services
After=mysql.service
Requires=mysql.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/sysconfig/bios
EnvironmentFile=-/etc/default/bios
ExecStart=@libexecdir@/bios/db-init

[Install]
WantedBy=multi-user.target
