# tools/bios-driver-netmon.service.in
[Unit]
Description=netmon driver for BIOS project
After=malamute.service
Requires=malamute.service

[Service]
Type=simple
ExecStart=@libdir@/bios/driver-netmon

[Install]
WantedBy=multi-user.target
