# tools/bios-agent-th.service.in
[Unit]
Description=Temperature and Humidity agent for BIOS project
After=malamute.service
Requires=malamute.service

[Service]
Type=simple
User=nobody
Group=dialout
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/sysconfig/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/sysconfig/bios__%n.conf
ExecStart=@libexecdir@/@PACKAGE@/agent-th

[Install]
WantedBy=multi-user.target
