# tools/bios-driver-nut.service.in
[Unit]
Description=NUT driver for BIOS project
After=malamute.service
Requires=malamute.service

[Service]
Type=simple
ExecStart=@libdir@/bios/driver-nut

[Install]
WantedBy=multi-user.target
