# tools/bios-agent-dbstore.service.in
[Unit]
Description=Database layer for BIOS
After=malamute.service bios-db-init.service
Requires=malamute.service bios-db-init.service

[Service]
Type=simple
User=nobody
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/sysconfig/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/sysconfig/bios__%n.conf
EnvironmentFile=-/etc/default/bios-db-rw
ExecStart=@libexecdir@/@PACKAGE@/agent-dbstore

[Install]
WantedBy=multi-user.target
