# tools/bios-agent-autoconfig.service.in
[Unit]
Description=Autoconfiguration agent for BIOS project
After=malamute.service bios-db-init.service
Requires=malamute.service
Requisite=bios-db-init.service

[Service]
Type=simple
ExecStart=@libexecdir@/bios/agent-autoconfig

[Install]
WantedBy=multi-user.target
