# tools/bios-server-agent.service.in
[Unit]
Description=Start BIOS agent that collects data about server it is running on
After=malamute.service
Requires=malamute.service

[Service]
Type=simple
ExecStart=@libexecdir@/bios/server-agent

