# tools/bios-agent-inventory.service.in
[Unit]
Description=Inventory agent for BIOS project
After=malamute.service bios-db-init.service
Requires=malamute.service bios-db-init.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/sysconfig/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/sysconfig/bios__%n.conf
ExecStart=@libexecdir@/bios/agent-inventory

[Install]
WantedBy=multi-user.target
