# tools/bios-networking.service.in
[Unit]
Description = Async restart of networking

[Service]
Type=simple
ExecStart=@libexecdir@/bios/bios-networking
