# tools/bios-agent-th.service.in
[Unit]
Description=Temperature and Humidity agent for BIOS project
After=malamute.service
Requires=malamute.service

[Service]
Type=simple
ExecStart=@libexecdir@/bios/agent-th

[Install]
WantedBy=multi-user.target
