# tools/bios-reset-button.service.in
[Unit]
Description=Handles events from reset button we have on our RC

[Service]
Type=simple
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/sysconfig/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/sysconfig/bios__%n.conf
ExecStart=@datadir@/@PACKAGE@/scripts/reset-button.sh

[Install]
WantedBy=bios.target
