# tools/biostimer-graphs-prefetch@.service.in
[Unit]
Description=BIOS-Timer service that prefetches average computations for quicker graphs with custom config file (implementation)
# This service only relies on REST API server directly
After=tntnet@bios.service bios-db-init.service
Requires=tntnet@bios.service bios-db-init.service

[Service]
Type=simple
ExecStart=@datadir@/@PACKAGE@/scripts/biostimer-graphs-prefetch.sh -C @sysconfdir@/@PACKAGE@/biostimer-graphs-prefetch@%I.conf

[Install]
WantedBy=multi-user.target
