# tools/bios-db.service.in
[Unit]
Description=Database layer for BIOS
After=malamute.service mysql.service
Requires=malamute.service mysql.service

[Service]
Type=simple
ExecStartPre=@libexecdir@/bios/db-init
ExecStart=@libexecdir@/bios/db-ng

