# tools/bios-agent-tpower.service.in
[Unit]
Description=Total power agent for BIOS project
After=malamute.service bios-db-init.service
Requires=malamute.service
Requisite=bios-db-init.service
# We only need the database during agetn-tpower startup to read
# our service configs, and don't care if it is restarted later

[Service]
Type=simple
ExecStart=@libexecdir@/bios/agent-tpower

[Install]
WantedBy=multi-user.target
