# tools/biostimer-outage.service.in
[Unit]
Description=BIOS-Timer service that re-counts outage statistics (implementation)
# This service only relies on database directly
After=bios-db-init.service
Requires=bios-db-init.service

[Service]
Type=simple
ExecStart=@libexecdir@/bios/outage

[Install]
WantedBy=multi-user.target
