# docs/examples/profile.d/bios_path.sh
# 42ity PATH
PATH="/usr/libexec/bios:$PATH"
export PATH
