# docs/examples/profile.d/fty_path.sh
# 42ity PATH
PATH="/usr/libexec/fty:$PATH"
export PATH
